# jscheme/src/build/crontab
45 3 * * * /bin/sh /home/groups/j/js/jscheme/htdocs/jscheme/src/build/update
